# recipes-example/mock/files/server.service
[Unit]
Description=simple mockup server

[Service]
ExecStart=/usr/bin/server

[Install]
WantedBy=multi-user.target